# mlops/crontab
2 23 * * * sh /workspaces/lgcns-mlops-practice/batch_prediction.sh
